# mmm_experiments/zmq-bridge.service
[Unit]
Description=0MQ bridge being used for MMM experiments

[Service]
User=xf28id1
Group=xf28id1
Environment="PATH=/home/xf28id1/project-mmm/mmm-experiments/venv/bin/:$PATH"
ExecStart=/home/xf28id1/project-mmm/mmm-experiments/venv/bin/python3 -m mmm_experiments.data.zmq_bridge
Restart=always

[Install]
WantedBy=multi-user.target
